# .config/systemd/user/autotrash.service
[Unit]
Description=Empty trash

[Service]
Type=oneshot
ExecStart="/usr/bin/autotrash" -d 8 -t
